# platform/nginx-java-21/docker/config/crontab
* * * * * cd /var/www && [program:command] >> /dev/stdout 2>&1

